# crontab
0 8 1 * * python /Users/anitaowens/Documents/GitHub/Python-Programming/choose_random_name.py
